# deploy/uwsgi.conf
# simple uWSGI script

description "uwsgi mednet-portal"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --ini /var/www/mednet-portal/deploy/uwsgi.ini
