# deploy/uwsgi.conf
# simple uWSGI script

description "uwsgi hamwan-portal"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --ini /var/www/hamwan-portal/deploy/uwsgi.ini
